# skills/macmini/install/com.macmini-skill.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macmini-skill.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>@@BINARY_PATH@@</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MACMINI_TOKEN_PATH</key>
        <string>@@TOKEN_PATH@@</string>
        <key>MACMINI_LOG_PATH</key>
        <string>@@LOG_PATH@@</string>
        <key>MACMINI_LISTEN_ADDR</key>
        <string>@@LISTEN_ADDR@@</string>
        <key>PATH</key>
        <string>@@PATH@@</string>
    </dict>

    <key>StandardOutPath</key>
    <string>@@LOG_PATH@@</string>

    <key>StandardErrorPath</key>
    <string>@@LOG_PATH@@</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
